# apps/api/deploy/supervisord.conf
; Supervisor config used in Fly.io single-machine deployment.
;
; Why one machine for two processes:
;   Fly's $5/month free credit covers ~1 shared-cpu-1x:512MB machine running
;   24/7. Splitting api and worker into separate machines doubles the bill.
;   Both processes have low CPU at MVP scale (a bar = a few RPS), so a single
;   container with two supervised processes is the right trade-off.
;
; Why supervisord and not "&" / tini multi-process / honcho:
;   - We need PID 1 to forward signals to children so Fly's graceful shutdown
;     gives uvicorn 15s to drain in-flight requests AND taskiq worker time to
;     ack the current task. supervisord handles SIGTERM correctly.
;   - We need to crash the whole machine (so Fly auto-restarts) if EITHER
;     process dies. Without that, a worker crash silently halts notifications
;     while the API keeps serving traffic.
;
; Logs go to stdout/stderr — Fly captures them and surfaces via `fly logs`.

[supervisord]
nodaemon=true
user=root
loglevel=info
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:api]
; Use absolute paths: supervisord does not run a login shell, and slim images can lack bash on PATH.
command=/usr/local/bin/uvicorn shiftops_api.main:app --host 0.0.0.0 --port 8000 --workers 1 --proxy-headers
directory=/app
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
exitcodes=0
priority=10

[program:worker]
command=/usr/local/bin/taskiq worker shiftops_api.infra.queue:broker --workers 1
directory=/app
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
exitcodes=0
priority=20
